# pkgbuild/auto-rules-io.install
post_install() {
    systemctl enable big-btrfs-balance.service
}

post_upgrade() {
    systemctl enable big-btrfs-balance.service
}

pre_remove() {
    systemctl disable big-btrfs-balance.service
} 
